# airootfs/etc/xdg/autostart/internet.desktop
[Desktop Entry]
Type=Application
Name=Refresh Pacman Databases
Exec=sudo /usr/bin/internet.sh
