# airootfs/etc/xdg/autostart/welcoming.desktop
[Desktop Entry]
Type=Application
Name=Reborn Welcome
Exec=gksudo antergos-welcome
Icon=gtk-info
